# tooling/builder/assets/clover2-firstboot.service
[Unit]
Description=Clover2 First Boot Script
After=network.target
Wants=network.target
After=cloud-init.target

[Service]
Type=oneshot
ExecStart=/bin/bash /root/clover2_firstboot.sh
RemainAfterExit=no
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
